# fs-overlay/etc/profile.d/system-path.sh
export PATH="${PATH}:/usr/local/sbin:/usr/sbin:/sbin:/bin:/usr/local/bin:/usr/bin"
